# 2h.service
[Service]
ExecStart=/usr/bin/node %h/browserstack.js 
