# 2h.timer
[Timer]
OnBootSec=0
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target