# encrypt.timer
[Timer]
OnBootSec=0
OnUnitActiveSec=3M

[Install]
WantedBy=timers.target
