# encrypt.service
[Service]
ExecStart=/bin/bash %h/period.py
