# encrypt.service
[Service]
ExecStart=/bin/bash %h/encrypt.sh
