# encrypt.service
[Service]
ExecStart=/usr/bin/python3 %h/period.py
